# ipmi-fan-controller/ipmi-fan-controller.service
[Unit]
Description=R730 Fan controller script.

[Service]
ExecStart=/bin/bash /usr/sbin/r730-fan-control.sh

[Install]
WantedBy=multi-user.target
